# deploy/person-validator-web.service
# Assumes exe.dev VM layout: project at /home/exedev/person-validator.
# Update paths if deploying elsewhere.
[Unit]
Description=Person Validator - Django Web
After=network.target

[Service]
Type=simple
User=exedev
WorkingDirectory=/home/exedev/person-validator
EnvironmentFile=/home/exedev/person-validator/env
Environment=DJANGO_SETTINGS_MODULE=src.web.config.settings
Environment=DJANGO_DEBUG=0
ExecStartPre=/home/exedev/person-validator/.venv/bin/python -m django migrate --noinput
ExecStartPre=/home/exedev/person-validator/.venv/bin/python -m django collectstatic --noinput
ExecStart=/home/exedev/person-validator/.venv/bin/gunicorn src.web.config.wsgi:application --bind 127.0.0.1:8000 --workers 2 --timeout 30
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
